# tools/maintenance/launchagents/com.brock.archive.compactor.monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>KEEP_DAYS</key>
		<string>90</string>
		<key>PURGE_OLD</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.brock.archive.compactor.monthly</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>/Users/brockmerkwan/Projects/devnotes/tools/maintenance/archive_compactor.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Day</key>
		<integer>1</integer>
		<key>Hour</key>
		<integer>19</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
</dict>
</plist>
